# eos-paygd/eos-paygd.service.in
[Unit]
Description=Pay As You Go Daemon
Documentation=man:eos-paygd(8)
# Only run on Phase 2 PAYG systems
ConditionKernelCommandLine=!eospayg
Before=systemd-user-sessions.service

[Service]
ExecStart=@libexecdir@/eos-paygd1
Type=dbus
BusName=com.endlessm.Payg1
NotifyAccess=main
# We must be root to avoid systemd's post-root-pivot killing spree
# and to make the eMMC boot partition writable
User=root
# Set KillMode and SendSIGKILL such that we're not stopped along with the associated target
KillMode=none
SendSIGKILL=no

# Sandboxing
# We might need access to /dev/mmcblkxboot0
DevicePolicy=auto
Environment=GIO_USE_VFS=local
Environment=GVFS_DISABLE_FUSE=1
Environment=GVFS_REMOTE_VOLUME_MONITOR_IGNORE=1
Environment=GSETTINGS_BACKEND=memory
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateNetwork=yes
PrivateTmp=yes
# We need PrivateUsers=no to be able to chown /var/lib/eos-payg
PrivateUsers=no
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=no
ProtectSystem=full
RestrictAddressFamilies=AF_UNIX
RestrictRealtime=yes
SystemCallErrorNumber=EPERM
SystemCallFilter=~@clock @cpu-emulation @debug @keyring @module @mount @obsolete @raw-io @resources
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
